# deploy/systemd/standx-openobserve-catchup.service
[Unit]
Description=StandX OpenObserve crash/boot catch-up (re-upload un-ingested log tail)
Documentation=file:/opt/standx/deploy/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/standx/maker.env
# When installed, the Stage 2 file overrides the same OpenObserve/log values so
# boot catch-up follows the A/B service's frozen deployment environment.
EnvironmentFile=-/etc/standx/maker-stage2-ab.env
WorkingDirectory=/opt/standx
# Idempotent: resumes each run's incremental checkpoint, so already-ingested
# events are not re-sent. The script always exits 0 (best-effort), so a down
# backend never blocks boot or the maker unit that orders after it.
ExecStart=/opt/standx/scripts/openobserve_catchup.sh

[Install]
WantedBy=multi-user.target
